# autobook.service
[Unit]
Description=AutoBook Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/AutoBook
EnvironmentFile=/root/AutoBook/.env
ExecStart=/root/AutoBook/.venv/bin/python -m bot.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
